# contrib/systemd/rosecd.service
[Unit]
Description=rosec - read-only Secret Service daemon
Documentation=https://github.com/jmylchreest/rosec
# Start after the D-Bus session bus is available.
# graphical-session.target ensures the Wayland compositor is running first,
# which is required for the unlock prompt GUI.
After=dbus.service graphical-session.target
Requires=dbus.service
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.secrets
# Adjust ExecStart to match your install location:
#   /usr/local/bin/rosecd   — system-wide install (sudo install)
#   %h/.cargo/bin/rosecd    — cargo install (development)
ExecStart=/usr/local/bin/rosecd
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
MemoryDenyWriteExecute=yes

# Allow network access (needed for Bitwarden API)
PrivateNetwork=no

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
